# decrypt-viewer/teslacam-decryptor.service
[Unit]
Description=Te_FITI-derived local decryptor + viewer for teslausb
After=network-online.target teslacam-latest-snapshot.service
Requires=teslacam-latest-snapshot.service

[Service]
Type=simple
Environment=DATA_DIR=/backingfiles/decrypt-viewer-state
ExecStartPre=/bin/mkdir -p /backingfiles/decrypt-viewer-state /backingfiles/decrypted
ExecStart=/usr/bin/python3 /opt/teslacam-decryptor/app/server.py \
  --scan /run/teslacam-latest/mnt/TeslaCam \
  --src /run/teslacam-latest/mnt/TeslaCam/EncryptedClips \
  --out /backingfiles/decrypted \
  --keys /backingfiles/decrypt-viewer-state/teslacam_keys.json \
  --port 8099 \
  --interval 60
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
