# decrypt-viewer/teslacam-sync.timer
[Unit]
Description=Run the NAS content sync every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Unit=teslacam-sync.service

[Install]
WantedBy=timers.target
