# decrypt-viewer/teslacam-sync.service
[Unit]
Description=Copy extra Tesla content (decrypted/LightShow/Wraps/Boombox/Music) to the NAS
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/teslacam-decryptor/sync-to-nas.sh
